# backend/systemd/detect-worker.service
[Unit]
Description=DETECT Bankr recipient intelligence worker
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/DETECT
EnvironmentFile=/opt/DETECT/.env
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/DETECT/.venv/bin/python -u /opt/DETECT/backend/src/main.py run-worker
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
